# deploy/systemd/mostify-stream-validator.timer
[Unit]
Description=Run the MOSTIFY strict stream validator every hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=120
Persistent=true
Unit=mostify-stream-validator.service

[Install]
WantedBy=timers.target
