# deploy/systemd/mostify-harvest-zambia.service
[Unit]
Description=MOSTIFY - fetch the broad Zambian radio catalog and import
After=network-online.target mostify-monitor.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/radio-monitor
EnvironmentFile=-/etc/mostify-monitor/mostify-monitor.env
ExecStart=/bin/bash -lc 'python3 scripts/zambia_station_harvest.py --max-probe 1400 --mytuner-pages 20 --include-no-icy --out scripts/data/zambia_harvest.json && /usr/bin/node scripts/import_zambia_stations.mjs --replace-stations-only scripts/data/zambia_harvest.json'
TimeoutStartSec=3600
# User=mostify
# Group=mostify

[Install]
WantedBy=multi-user.target
