# deploy/systemd/mostify-orb-poller.timer
[Unit]
Description=Run MOSTIFY ORB track poller every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
