# deploy/systemd/mostify-orb-poller.service
[Unit]
Description=MOSTIFY OnlineRadioBox track poller (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/radio-monitor
EnvironmentFile=-/etc/mostify-monitor/mostify-monitor.env
# Replace with: which node
ExecStart=/usr/bin/node scripts/orb_track_poller.mjs
# User=mostify
# Group=mostify

[Install]
WantedBy=multi-user.target
