# deploy/systemd/mostify-harvest-zambia.service
[Unit]
Description=MOSTIFY — fetch every Zambian radio station from radio-browser and import
After=network-online.target mostify-monitor.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/radio-monitor
EnvironmentFile=-/etc/mostify-monitor/mostify-monitor.env
ExecStart=/usr/bin/node scripts/zambia_radio_browser_harvest.mjs --import --validate
TimeoutStartSec=2700
# User=mostify
# Group=mostify

[Install]
WantedBy=multi-user.target
